# deploy/systemd/stackchan-touch-responder.service
[Unit]
Description=StackChan Immediate Touch Responder
After=network-online.target mosquitto.service
Wants=network-online.target

[Service]
Type=simple
User=stackchan
Group=stackchan
WorkingDirectory=/opt/stackchan-remote-mcp/server
EnvironmentFile=/etc/stackchan-remote-mcp.env
ExecStart=/opt/stackchan-remote-mcp/.venv/bin/python /opt/stackchan-remote-mcp/server/touch_responder.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
