# deploy/systemd/stackchan-relay.service
[Unit]
Description=StackChan Photo Relay
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=stackchan
Group=stackchan
WorkingDirectory=/opt/stackchan-remote-mcp/server
EnvironmentFile=/etc/stackchan-remote-mcp.env
ExecStart=/opt/stackchan-remote-mcp/.venv/bin/python /opt/stackchan-remote-mcp/server/photo_relay.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
